# wallpaper-update.timer
[Unit]
Description=Wallpaper update

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
