# wallpaper-update.timer
[Unit]
Description=Wallpaper update

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
